# application/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stopsignal=QUIT

[program:horizon]
command=php /var/www/html/artisan horizon
autostart=true
autorestart=true
